# sdcard/lab/android/etc/crontabs/root
30 11,19 * * * /data/xbin/daily-jobs

# * * * * * bash /sdcard/lab/android/bin-extra/debug &> /sdcard/lab/android/log/debug.log